# apps/core/agenda/setup/etc/cron.daily/elastix_contacs_cleanup
#!/bin/sh
find /var/www/elastixdir/contacts_images/ -type f -name "tmp_contact_*" -exec rm -rf {} \; >/dev/null 2>&1
