# apps/core/pbx/setup/etc/cron.daily/asterisk_cleanup
#!/bin/sh
/usr/bin/find /var/lib/asterisk/sounds/tts -type f -mtime +1 | /usr/bin/xargs /bin/rm -f >/dev/null 2>&1
